# root/etc/s6-overlay/s6-rc.d/rdp-cuse/run
#!/bin/bash

[ -e /config/logs/ ] || mkdir /config/logs

while true; do
  [ -e /dev/cuse ] && {
    /root/cuse/scripts/localstart.sh >> /config/logs/cuse.run 2>&1
  }
  sleep 1
done

