# containerfiles/code/root/etc/s6-overlay/s6-rc.d/nginx-proxy/run
#!/bin/bash

[ -e /config/logs/ ] || mkdir /config/logs

/usr/sbin/nginx >/config/logs/00_nginx.out 2>/config/logs/00_nginx.err

